# etc/services.d/dovecot/run
#!/bin/bash

cat /private/ssl/{cert,ca}.pem > /etc/dovecot/bundle.pem
exec dovecot -F
